# tools/wlconfignew
#!/bin/sh
### BEGIN INIT INFO
# Provides:          wlconfignew
# Required-Start:
# Required-Stop:
# Default-Start:
# Default-Stop:
# Short-Description: Configure a Broadcom "wl" WLAN interface from its wl.conf.<iface> (ssid/method/key)
### END INIT INFO

set -eu

PATH="/sbin:/bin:/usr/sbin:/usr/bin"

WL="/usr/bin/wl"
IFACE=""
CONF=""

while getopts "i:c:" opt; do
  case "$opt" in
    i) IFACE="$OPTARG" ;;
    c) CONF="$OPTARG" ;;
    *)
      echo "Usage: $0 -i<iface> -c<wl.conf.iface>" >&2
      exit 2
      ;;
  esac
done

case "$IFACE" in
  wlan[0-9]*) ;;
  *)
    echo "Usage: $0 -i<iface> -c<wl.conf.iface>" >&2
    exit 2
    ;;
esac

if [ -z "$CONF" ] || [ ! -f "$CONF" ]; then
  echo "No wl config found at '$CONF'" >&2
  exit 1
fi

# wl.conf.<iface> is a plain "key=value" file with exactly one network:
#   ssid=...
#   method=...   (None, wep, wpa, wpa2)
#   key=...
METHOD="None"
SSID=""
KEY=""

while IFS='=' read -r key value; do
  case "$key" in
    ssid) SSID="$value" ;;
    method) METHOD="$value" ;;
    key) KEY="$value" ;;
  esac
done < "$CONF"

if [ -z "$SSID" ]; then
  echo "No ssid set in $CONF" >&2
  exit 1
fi

# ifconfig up already ran as the previous pre-up command; give the driver the
# same settle time the original vendor wl-config.sh relied on.
sleep 0.1

"$WL" wowl_clear
"$WL" up
"$WL" auth 2
"$WL" infra 1

case "$METHOD" in
  wep)
    "$WL" wsec 1
    "$WL" wpa_auth 0
    "$WL" addwep 0 "$KEY"
    ;;
  wpa)
    "$WL" wsec 6
    "$WL" wpa_auth 4
    "$WL" sup_wpa 1
    "$WL" set_pmk "$KEY"
    ;;
  wpa2)
    "$WL" wsec 6
    "$WL" wpa_auth 128
    "$WL" sup_wpa 1
    "$WL" set_pmk "$KEY"
    ;;
  *)
    "$WL" wsec 0
    "$WL" wpa_auth 0
    ;;
esac

"$WL" ssid "$SSID"

count=1
while [ "$count" -le 10 ]; do
  count=$((count + 1))
  status="$("$WL" status 2>/dev/null | cut -c1-4)"
  [ "$status" = "SSID" ] && break
  sleep 0.5
done

exit 0
